# deploy/launchd/com.tailhub.hub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Tailhub launchd agent. Installed by scripts/install-hub-startup.sh into
     ~/Library/LaunchAgents/com.tailhub.hub.plist — the installer replaces
     @NODE@, @CLI@, and @LOGDIR@ with resolved absolute paths. @LOGDIR@ is a
     per-user 0700 directory: the hub prints its generated admin token to
     stdout on first start, so logs must not land in a world-readable path. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tailhub.hub</string>
  <key>ProgramArguments</key>
  <array>
    <string>@NODE@</string>
    <string>@CLI@</string>
    <string>start</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>StandardOutPath</key>
  <string>@LOGDIR@/tailhub.out.log</string>
  <key>StandardErrorPath</key>
  <string>@LOGDIR@/tailhub.err.log</string>
</dict>
</plist>
